# ds1307/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable ds1307sync
systemctl restart ds1307sync